# services/text-to-speech/text-to-speech.service
[Unit]
Description=Text-to-Speech Service (Piper + Edge TTS)
After=network.target

[Service]
Type=simple
User=aicrm
WorkingDirectory=/opt/ai-assistant/services/assistant-core/services/text-to-speech
ExecStart=/opt/ai-assistant/services/assistant-core/services/text-to-speech/venv/bin/uvicorn app:app --host 0.0.0.0 --port 5400
Restart=always
RestartSec=5
Environment=PIPER_DIR=/opt/piper/piper
Environment=PIPER_MODELS_DIR=/opt/piper/models
Environment=PIPER_DEFAULT_VOICE=en_US-lessac-medium

[Install]
WantedBy=multi-user.target
